# pmaports/main/postmarketos-bootsplash/postmarketos-bootsplash-openrc.post-install
#!/bin/sh

rc-update add kill-plymouth default
